# packaging/debian/postinst
#!/bin/bash
set -e

if [ "$1" = "configure" ]; then
    # Create ulog system user if not exists
    if ! getent passwd ulog >/dev/null; then
        adduser --system --group --no-create-home --shell /usr/sbin/nologin ulog
    fi

    # Add ulog to serial port groups (varies by distro)
    # Arch: uucp, Debian/Ubuntu: dialout
    for group in uucp dialout tty; do
        if getent group "$group" >/dev/null; then
            usermod -aG "$group" ulog
        fi
    done

    # Set config file permissions (root:ulog 0640)
    chown root:ulog /etc/ulog.conf
    chmod 0640 /etc/ulog.conf

    # Create state directory for session tracking
    mkdir -p /var/lib/ulog/sessions
    chown ulog:ulog /var/lib/ulog /var/lib/ulog/sessions
    chmod 0750 /var/lib/ulog /var/lib/ulog/sessions

    # Generate initial config (creates service file and udev rule)
    /usr/bin/ulog-genconfig

    systemctl daemon-reload
    udevadm control --reload-rules

    echo ""
    echo "ulog installed."
    echo "Enable with: systemctl enable --now ulog.service ulog-rollover.timer ulog-genconfig.path"
    echo "Config: /etc/ulog.conf (root:ulog 0640)"
    echo ""
fi

#DEBHELPER#
