# services/ulog-rollover.service
[Unit]
Description=Restart ulog for day rollover

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart ulog.service
