# services/ulog-rollover.timer
[Unit]
Description=Restart ulog at midnight for day rollover

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
